# tools/taskwarrior/systemd/timewsync.timer
[Unit]
Description=Run timewsync periodically and on boot

[Timer]
OnBootSec=15m
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target

