# tools/taskwarrior/systemd/timewsync.service
[Unit]
Description=timewsync

[Service]
ExecStart=/home/sebas/.asdf/shims/timewsync
Type=oneshot

[Install]
WantedBy=default.target
