# tools/taskwarrior/systemd/timewsync.service
[Unit]
Description=timewsync

[Service]
ExecStart="%h/.asdf/shims/timewsync"
Type=oneshot

[Install]
WantedBy=default.target
