# tools/taskwarrior/systemd/tw-sync.service
[Unit]
Description=tw-sync

[Service]
ExecStart="%h/.local/bin/sync-tw"
Type=oneshot

[Install]
WantedBy=default.target
