# tools/taskwarrior/systemd/tw-sync.service
[Unit]
Description=tw-sync

[Service]
ExecStart=/usr/local/bin/task sync
Type=oneshot

[Install]
WantedBy=default.target
